# dot_config/systemd/user/restic-mc-instance-backup.timer
[Unit]
Description=Offsite MC instance backup to restic

[Timer]
OnCalendar=*-*-* *:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
